# www-apps/jackett-bin/files/jackett.init
#!/sbin/openrc-run

pidfile="/run/jackett.pid"
command="/opt/jackett/jackett"
command_args="--ListenPrivate --NoUpdates"
command_user="jackett:jackett"
command_background=true
