# net-p2p/freenet/files/freenet
#!/sbin/openrc-run
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start() {
        su freenet --shell /bin/sh -c "/var/freenet/run.sh start"
	eend $?
}

status() {
        su freenet --shell /bin/sh -c "/var/freenet/run.sh status"
	eend $?
}

stop() {
        su freenet --shell /bin/sh -c "/var/freenet/run.sh stop"
	eend $?
}
